# dracut/28immucore/immucore.rc
#!/sbin/openrc-run

name="immucore"
command="/usr/bin/immucore"

depend() {
        need udev
}
